# systemd/sync-nodes.service
[Unit]
Description=Sync Xboard nodes to local systemd instances

[Service]
Type=oneshot
ExecStart=/usr/local/bin/sync-nodes
