# systemd/update-xboard-node.service
[Unit]
Description=Update Xboard-Node automatic

[Service]
Type=oneshot
ExecStart=/usr/local/bin/update-xboard-node.sh
